# runit/fetchcord/run
#!/bin/sh
exec python3 -u -m fetch_cord.run-rpc
